# git-gc-chromium-src.service
[Unit]
Description="git gc chromium"

[Service]
ExecStart=/usr/bin/git -C /media/src/chromium/src gc
