# git-gc-depot-tools.service
[Unit]
Description="git gc depot_tools"

[Service]
ExecStart=/usr/bin/git -C %h/src/headless-shell/out/depot_tools gc --aggressive
