# headless-shell.service
[Unit]
Description="build headless-shell"

[Service]
ExecStart=/home/user/src/headless-shell/build.sh -j 14 -p
