# headless-shell.service
[Unit]
Description="headless-shell build"

[Service]
ExecStart=/home/user/src/headless-shell/build.sh -j 14 -u
