# podman-image-prune.service
[Unit]
Description="prune podman images"

[Service]
ExecStart=/usr/bin/podman image prune --force
